# docker/etc/supervisor.d/web_server.conf
[program:nginx]
command=/etc/init.d/nginx start
autorestart=true
user=nginx
